# systemd/notify_touchy.service
[Unit]
# Not finished. idk

[Service]
ExecStart=/usr/local/bin/notify_touchy

[Install]
WantedBy=default.target
